# root/etc/services.d/sshvpn/run
#!/usr/bin/with-contenv bash

openvpn /config/vpn.ovpn
